# opscenter/src/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/lib/supervisord.log

[program:opscenter]
user=root
command=/usr/share/opscenter/bin/opscenter
