# opscenter/src/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/log/supervisord.pid

[program:opscenter]
user=root
command=/usr/share/opscenter/bin/opscenter
